# linux/domofox.service
[Unit]
Description=Domofox API server
After=network.target

[Service]
Type=exec
WorkingDirectory=/opt/domofox/api
ExecStart=/usr/bin/node /opt/domofox/api/index.js
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
